# ops/mlx/com.dagagent.mlx-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dagagent.mlx-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/andrew/.venvs/rapid-mlx/bin/rapid-mlx</string>
        <string>serve</string>
        <string>mlx-community/Qwen3.6-35B-A3B-4bit</string>
        <string>--served-model-name</string>
        <string>qwen3.6-35b-moe</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
        <string>--max-tokens</string>
        <string>12000</string>
        <string>--max-num-seqs</string>
        <string>1</string>
        <string>--gpu-memory-utilization</string>
        <string>0.72</string>
        <string>--kv-cache-quantization</string>
        <string>--kv-cache-quantization-bits</string>
        <string>4</string>
        <string>--prefill-step-size</string>
        <string>8192</string>
        <string>--pin-system-prompt</string>
        <string>--reasoning-parser</string>
        <string>qwen3</string>
        <string>--enable-auto-tool-choice</string>
        <string>--tool-call-parser</string>
        <string>auto</string>
        <string>--default-temperature</string>
        <string>0.1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrew</string>

    <!-- Restart automatically on crash -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 30s between restarts to avoid rapid crash loops -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/andrew/Library/Logs/dagagent-mlx-server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrew/Library/Logs/dagagent-mlx-server.err.log</string>
</dict>
</plist>
